# content/marketing/routines/launchd/ai.thepsychology.linkedin-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.thepsychology.linkedin-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/anderschan/.thepsychology-automation/run-linkedin-daily.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/anderschan/.thepsychology-automation/linkedin-daily.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/anderschan/.thepsychology-automation/linkedin-daily.launchd.err</string>
</dict>
</plist>
